# crontab
55 * * * * . /check_for_data.sh        > /var/log/cron_check_for_data.log 2>&1
